# gpio-control.service
[Unit]
Description=GPIO Control Server for Door Lock
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 /home/pi/resonance-password-screen/gpio-server.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
